# tools/etc/rc.local
#!/bin/sh
screen -Dm /bin/openwebrx.sh &
